# .docker/supervisor/000-queue-iopro.conf
[program:queue-iopro]
autostart=true
autorestart=true
numprocs=5
process_name=%(process_num)s
command=php /var/www/public/index.php queue doctrine iopro --start
